# restart-nightly.service
[Service]
ExecStart=%h/unique-kit/restart-when-empty
