# update-maps.service
[Service]
ExecStart=%h/unique-kit/update-maps
